# .build/pianocktail.deb-1.0/debian/postrm
#!/bin/sh
if [ "$1" = purge ]; then
        update-rc.d pianocktail.deb remove >/dev/null
fi
