# .build/pianocktail.deb-1.0/debian/postinst
#!/bin/sh
update-rc.d pianocktail.deb defaults >/dev/null
invoke-rc.d pianocktail.deb start
