# .build/pianocktail.deb-1.0/debian/prerm
#!/bin/sh
invoke-rc.d pianocktail.deb stop
