# .build/twisted-pyanocktaild-1.0/debian/postinst
#!/bin/sh
update-rc.d twisted-pyanocktaild defaults >/dev/null
invoke-rc.d twisted-pyanocktaild start
