# .build/twisted-pyanocktaild-1.0/debian/postrm
#!/bin/sh
if [ "$1" = purge ]; then
        update-rc.d twisted-pyanocktaild remove >/dev/null
fi
