# .build/twisted-pyanocktaild-1.0/debian/prerm
#!/bin/sh
invoke-rc.d twisted-pyanocktaild stop
